# roles/mysql/templates/mysql.backup.cron.j2
{% if inventory_hostname == mysql_backup_host %}
50 21 * * * backup mysqldump agama > /home/backup/mysql/agama.sql
55 22 * * * backup duplicity --no-encryption full /home/backup/mysql/ rsync://davidayson8@backup/mysql
{% endif %}
